# security-monitoring/network-watchdog.service
[Unit]
Description=LoxProx Network Stack Watchdog
Documentation=https://github.com/sgtsilver/loxprox
After=network.target nginx.service

# This is a systemd SYSTEM service (not --user). It runs as root by default,
# which is required to restart networking.service and call reboot.
# No sudo configuration is needed. No passwordless access is granted.

# The watchdog timer fires every 60 s; StartLimitBurst=0 disables start-rate
# limiting so systemd never refuses to run a cycle. There is intentionally NO
# FailureAction (see [Service]) — nothing here reboots on a failed start.
StartLimitIntervalSec=600
StartLimitBurst=0

[Service]
Type=oneshot
ExecStart=/opt/loxprox/network-watchdog.sh

# H3: FailureAction=reboot was REMOVED. systemd fires FailureAction immediately on
# entering the failed state (StartLimitBurst gates StartLimitAction, NOT
# FailureAction), so EVERY non-zero exit rebooted — a corrupt state file, a full
# disk breaking log writes, or the old give-up `exit 1` all became reboot loops.
# The script owns its own gated reboot (handle_reboot, max 2/hour) and now exits 0
# on give-up. A broken script must alert/log, never silently reboot the box.

# Logging
StandardOutput=append:/var/log/loxprox-network-watchdog.log
StandardError=append:/var/log/loxprox-network-watchdog.log
SyslogIdentifier=loxprox-watchdog
