# security-monitoring/network-watchdog.service
[Unit]
Description=LoxProx Network Stack Watchdog
Documentation=https://github.com/sgtsilver/loxprox
After=network.target nginx.service

# This is a systemd SYSTEM service (not --user). It runs as root by default,
# which is required to restart networking.service and call reboot.
# No sudo configuration is needed. No passwordless access is granted.

# systemd-level safety net: if the script exits non-zero, systemd can
# trigger a reboot as last-resort recovery. StartLimitBurst=0 is required
# because the timer fires every 60 s; a finite burst limit would block the
# 4th start (~3 min) and falsely trigger FailureAction=reboot. The script
# already has its own anti-loop counter (max 2 reboots/hour). Belongs in
# [Unit] per systemd.unit(5); placing them in [Service] is silently ignored.
StartLimitIntervalSec=600
StartLimitBurst=0

[Service]
Type=oneshot
ExecStart=/opt/loxprox/network-watchdog.sh

# If the script exits non-zero after the burst limit is exhausted,
# systemd triggers a reboot as a last-resort recovery measure.
FailureAction=reboot

# Logging
StandardOutput=append:/var/log/loxprox-network-watchdog.log
StandardError=append:/var/log/loxprox-network-watchdog.log
SyslogIdentifier=loxprox-watchdog
